# deploy/sendguard-agent.service
[Unit]
Description=SendGuard Agent — protección SMTP/IMAP para Zimbra
Documentation=https://github.com/perulinux/sendguard
After=network.target zimbra.service
Wants=zimbra.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/sendguard-agent -config /etc/sendguard/agent.yaml
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sendguard-agent

# Seguridad: limitar superficie de ataque del proceso
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/zimbra/common/conf /var/log /var/lib/sendguard /etc/ufw
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
